# com.grantwriterpack.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS LaunchAgent: daily grants.gov scan at 7:00 AM.

  To install:
    cp com.grantwriterpack.daily.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.grantwriterpack.daily.plist

  To uninstall:
    launchctl unload ~/Library/LaunchAgents/com.grantwriterpack.daily.plist
    rm ~/Library/LaunchAgents/com.grantwriterpack.daily.plist

  To run manually now:
    launchctl start com.grantwriterpack.daily
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.grantwriterpack.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>INSTALL_DIR/scripts/daily-scan.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>INSTALL_DIR/data/daily-scan.log</string>

    <key>StandardErrorPath</key>
    <string>INSTALL_DIR/data/daily-scan.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
